# scripts/julliuz-bot.service
[Unit]
Description=Julliuz Bot Service
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/julliuz-bot
Environment="PATH=/root/julliuz-bot/venv/bin"
Environment="PYTHONPATH=/root/julliuz-bot"
ExecStart=/root/julliuz-bot/venv/bin/python /root/julliuz-bot/scripts/run_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target